# files/beszel-agent.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
PROG=/usr/bin/beszel-agent

append_fs() {
    local value="$1"
    if [ -z "$FS_LIST" ]; then
        FS_LIST="$value"
    else
        FS_LIST="$FS_LIST,$value"
    fi
}

append_nic() {
    local value="$1"
    if [ -z "$NIC_LIST" ]; then
        NIC_LIST="$value"
    else
        NIC_LIST="$NIC_LIST,$value"
    fi
}

start_service() {
    local enable hub_url key listen log_level port token skip_gpu skip_systemd system_name

    config_load beszel-agent

    config_get_bool enable "main" "enable" 0
    [ "$enable" -eq 1 ] || return 0

    config_get listen "main" "listen" ""
    config_get log_level "main" "log_level" "info"
    config_get port "main" "port" "45876"
    config_get key "main" "key" ""
    config_get token "main" "token" ""
    config_get skip_gpu "main" "skip_gpu" "false"
    config_get skip_systemd "main" "skip_systemd" "false"
    config_get hub_url "main" "hub_url" ""
    config_get system_name "main" "system_name" ""

    FS_LIST=""
    config_list_foreach "main" "extra_filesystems" append_fs

    NIC_LIST=""
    config_list_foreach "main" "nics" append_nic

    procd_open_instance "beszel-agent"
    procd_set_param command "$PROG"

    procd_set_param env HOST_PROC="/proc" HOST_SYS="/sys"

    if [ -n "$listen" ]; then
        procd_append_param env LISTEN="$listen"
    else
        procd_append_param env PORT="$port"
    fi

    [ -n "$system_name" ] && procd_append_param env SYSTEM_NAME="$system_name"
    [ -n "$key" ] && procd_append_param env KEY="$key"
    [ -n "$token" ] && procd_append_param env TOKEN="$token"
    [ -n "$hub_url" ] && procd_append_param env HUB_URL="$hub_url"
    [ -n "$log_level" ] && procd_append_param env LOG_LEVEL="$log_level"
    [ -n "$skip_gpu" ] && procd_append_param env SKIP_GPU="$skip_gpu"
    [ -n "$skip_systemd" ] && procd_append_param env SKIP_SYSTEMD="$skip_systemd"
    [ -n "$FS_LIST" ] && procd_append_param env EXTRA_FILESYSTEMS="$FS_LIST"
    [ -n "$NIC_LIST" ] && procd_append_param env NICS="$NIC_LIST"
    
    procd_set_param stdout 1
    procd_set_param stderr 1

    procd_set_param respawn
    procd_close_instance
}

service_triggers() {
    procd_add_reload_trigger "beszel-agent"
}

reload_service() {
    stop
    start
}
